# backup.timer
[Unit]
Description=Run backup after boot and then on a weekly bases

[Timer]
OnBootSec=10min
OnUnitActiveSec=1w

[Install]
WantedBy=backup.target
